# sunshine.service
#Lives in ~/.config/systemd/user
[Unit]
Description=Syncthing

[Service]
ExecStart=/usr/bin/sunshine

[Install]
WantedBy=default.target
